# docs/nginx_custom_domains.conf
# Nginx Configuration for CNAME-based SaaS Tracking

To support 1000+ custom domains pointing to your Laravel backend, use the following configuration.

## 1. Main Application & Dynamic Domain Server
Add this server block to your Nginx configuration (e.g., `/etc/nginx/sites-available/tracking.conf`).

```nginx
server {
    listen 80;
    listen [::]:80;
    
    # Catch-all for all custom domains pointing to this server
    server_name _; 

    root /var/www/html/public;
    index index.php index.html;

    # Handle SSL via Certbot (Certbot will add its own directives here)
    
    location / {
        try_files $uri $uri/ /index.php?$query_string;
    }

    location ~ \.php$ {
        include snippets/fastcgi-php.conf;
        fastcgi_pass unix:/var/run/php/php8.2-fpm.sock;
        fastcgi_param SCRIPT_FILENAME $document_root$fastcgi_script_name;
        include fastcgi_params;
    }

    location ~ /\.ht {
        deny all;
    }
}
```

## 2. Automated SSL Provisioning
For a SaaS with many custom domains, use **Certbot with Nginx** in a scriptable way.

### Option A: Certbot (Standard)
When a user clicks "Connect Domain", you can trigger a background job to run:
```bash
sudo certbot --nginx -d data.user-site.com --non-interactive --agree-tos --email your-admin@email.com
```

### Option B: Wildcard (If using subdomains of your own domain)
If you are only allowing `*.yourdomain.com`, use a wildcard SSL.

---

# Maintenance & Scaling Guide

## 1. Supervisor for Queue Workers
Ensure your queue worker is always running to process events asynchronously.
Create `/etc/supervisor/conf.d/laravel-worker.conf`:

```ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
stopwaitsecs=3600
```

## 2. Redis Configuration
Ensure Redis is configured to handle high throughput. Increase the `maxclients` if necessary in `redis.conf`.

## 3. Rate Limiting
The system already uses Redis-based rate limiting per project. Monitor Redis memory usage to ensure limits are tracked correctly.
